# atem-html-source.service
[Unit]
Description=ATEM HTML Source Server
After=network.target

[Service]
Type=simple
User=USERNAME
WorkingDirectory=REPO_DIR
ExecStart=/usr/bin/env node server.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production PORT=3000

[Install]
WantedBy=multi-user.target
